# crontab/user
#######################
# Update configs repo #
#######################

# crontab
3 */2 * * * cd ~/Desktop/configs/crontab && crontab -l > user

###########################
# Refresh direnv database #
###########################
@reboot direnv prune

###############
# Empty trash #
###############
@reboot osascript -e 'tell application "Finder" to empty trash'

########
# brew #
########
@reboot brew autoremove
@reboot brew cleanup

########
# tldr #
########
@reboot tldr --update

